# supervisor.ini
[program:jc_flask]
command=uwsgi --socket 127.0.0.1:3031 --file runserver.py --callable app --processes 1 --venv /Users/jcurle/.virtualenvs/jamiecurle/ 
directory=/Users/jcurle/Sites/jamiecurle
user=jcurle
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT