# supervisor.ini
[program:jc_flask]
command=/home/curle/.virtualenvs/jamiecurle/bin/uwsgi --http 127.0.0.1:26331 --file runserver.py --callable app --processes 1 --venv /home/curle/.virtualenvs/jamiecurle/ 
directory=/home/curle/sites/jamiecurle/jamiecurle
user=jcurle
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT